# services/sddm.service.in
[Unit]
Description=Simple Desktop Display Manager
Documentation=man:sddm(1) man:sddm.conf(5)
Conflicts=getty@tty1.service
After=systemd-user-sessions.service getty@tty1.service plymouth-quit.service

[Service]
ExecStart=@CMAKE_INSTALL_FULL_BINDIR@/sddm
Restart=always

[Install]
Alias=display-manager.service

