# services/sddm.service.in
[Unit]
Description=Simple Desktop Display Manager
Documentation=man:sddm(1) man:sddm.conf(5)
Conflicts=getty@tty${SDDM_INITIAL_VT}.service
After=systemd-user-sessions.service getty@tty${SDDM_INITIAL_VT}.service plymouth-quit.service systemd-logind.service
PartOf=graphical.target
StartLimitIntervalSec=30
StartLimitBurst=2

[Service]
ExecStart=@CMAKE_INSTALL_FULL_BINDIR@/sddm
Restart=always

[Install]
Alias=display-manager.service

